# scripts/98-moni.rules
ACTION=="add", SUBSYSTEM=="drm", KERNEL=="card0", RUN+="/home/phon/home/phonkd/git/dotfiles/scripts/amonitor.sh"

